# cornguard.service
[Unit]
Description=CornGuard FastAPI (TensorFlow CPU, gunicorn+uvicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/site1/backend
Environment="PATH=/var/www/site1/backend/.venv/bin"
Environment="OMP_NUM_THREADS=1"
Environment="TF_NUM_INTRAOP_THREADS=1"
Environment="TF_NUM_INTEROP_THREADS=1"
Environment="TF_CPP_MIN_LOG_LEVEL=2"
ExecStart=/var/www/site1/backend/.venv/bin/gunicorn main:app \
    --worker-class uvicorn.workers.UvicornWorker \
    --workers 1 \
    --bind 0.0.0.0:8001 \
    --timeout 120 \
    --graceful-timeout 30
Restart=always
RestartSec=5
LimitNOFILE=65535
# Optional resource guards:
# MemoryMax=1G
# CPUQuota=80%

[Install]
WantedBy=multi-user.target
